# firstboot.service
[Unit]
Description=Perform initial boot setup
[Service]
Type=simple
ExecStart=/firstboot.sh
[Install]
WantedBy=multi-user.target
